# modules/base/files/cron/puppet.cron
#!/usr/bin/perl -w 
use strict;
use warnings;

if(`pgrep -u root ruby-mri | wc -l` != 1) {
	system("puppet agent --server puppet");
} 
